# automathemely/lib/installation_files/sun-times.timer
[Unit]
Description=Update automathemely sun times daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
